# toolboxes/bluefin-cli/files/etc/profile.d/homebrew-env.sh
# shellcheck shell=sh

PATH="$PATH:/home/linuxbrew/.linuxbrew/sbin"
export PATH
HOMEBREW_REPOSITORY="/home/linuxbrew/.linuxbrew"
export HOMEBREW_REPOSITORY
HOMEBREW_PREFIX="/home/linuxbrew/.linuxbrew"
export HOMEBREW_PREFIX
HOMEBREW_CELLAR="/home/linuxbrew/.linuxbrew/Cellar"
export HOMEBREW_CELLAR