# toolboxes/bluefin-cli/files/etc/profile.d/zz-profile_sourced.sh
# shellcheck shell=sh
if [ -z "$PROFILESOURCED" ] && [ "$PS1" ]; then
  PROFILESOURCED="Y"
fi
